# phase-a/systemd/govless-bot.service
[Unit]
Description=goVLESS Telegram bot (aiogram)
Documentation=https://github.com/anten-ka/govless
After=govlessctl.service network-online.target
Wants=network-online.target
Requires=govlessctl.service

[Service]
Type=simple
User=govless
Group=govless
EnvironmentFile=/etc/govless/bot.env
ExecStart=/usr/local/bin/govless-bot
Restart=on-failure
RestartSec=10

# Hardening (014 Rule 3a)
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
ReadOnlyPaths=/etc/govless
ReadWritePaths=/var/lib/govless /run/govless
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictRealtime=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
