# phase-a/systemd/tunnel-health.service
[Unit]
Description=goVLESS Quick Tunnel health check (hysteresis + cold-start grace)
After=cloudflared-quick.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/govless-tunnel-health-check
StandardOutput=journal
StandardError=journal
